# source/bin/zwt
#!/bin/sh
### BEGIN INIT INFO
#	Provides:		zwt
#	Required-Start:		$all
#	Required-Stop:		
#	Default-Start:		2 3 4 5
#	Default-Stop:		0 1 6
#	Short-Description:	Zero Wine Tryouts
### END INIT INFO

ZWT_NAME=zerowine
BOA_NAME=boa
USER=malware
GROUP=malware
SAMPLE_DIR=/tmp/vir

ZWT_PATH=/home/$USER/$ZWT_NAME/
BOA_PATH=/usr/bin/$BOA_NAME
BOA_ARGS=" -c $ZWT_PATH"

[ -x "$BOA_PATH" ] || exit 0
[ -x "$ZWT_PATH" ] || exit 0

set -e

case "$1" in
	start)
		echo -n "Starting Zero Wine Tryouts..."
		su $USER -c "mkdir -p $SAMPLE_DIR"
		start-stop-daemon --start --quiet \
		 --chuid $USER:$GROUP --umask 0022 \
		 --name $BOA_NAME --exec $BOA_PATH -- $BOA_ARGS
		echo "."
		;;
	stop)
		echo -n "Stopping Zero Wine Tryouts..."
		start-stop-daemon --stop --quiet \
		 --name $BOA_NAME
		echo "."
		;;
	*)
		echo "Usage: /etc/init.d/$0 {start|stop}"
		exit 1
esac

exit 0
